# livekit-agent-service/supervisord.conf
; LiveKit Agent Service — roda dois processos:
;   - worker  : long-lived agent worker registrado no LiveKit Cloud
;   - dispatcher : Flask HTTP API que cria salas + dispara SIP outbound
;
; supervisord coloca os dois sob mesmo lifecycle (se um morrer, container reinicia)

[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0

[program:worker]
command=python -m livekit.agents.cli connect --url %(ENV_LIVEKIT_WS_URL)s worker.py
directory=/app
autostart=true
autorestart=true
startretries=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED=1

[program:dispatcher]
command=python -m flask --app dispatcher run --host=0.0.0.0 --port=5042
directory=/app
autostart=true
autorestart=true
startretries=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED=1,FLASK_APP=dispatcher
